# custom/randomdis
#!/sbin/openrc-run

# Copyright (c) Natanael Copa
# This code is licensed under BSD-2-Clause

description="custom script modification"


start() {
		echo 0 > /proc/sys/kernel/randomize_va_space
		echo 0 > /proc/sys/kernel/yama/ptrace_scope 
		#intel_gpu_frequency -c min=500
		for i in /sys/devices/system/cpu/cpufreq/policy*;do
			# echo "schedutil" > $i/scaling_governor
			echo $(cat /sys/devices/system/cpu/cpufreq/policy0/base_frequency) > $i/scaling_max_freq
			# echo 2000000 > $i/scaling_max_freq
		done
	eend $?
}
